# source/usr/local/emhttp/plugins/mirsht/defaults/default.cron
# MirrorShuttle:
* * * * * /usr/bin/mirrorshuttle-cron > /dev/null 2>&1
